# templates/default/systemd-dashboard.service.erb
# Generated by Chef - DO NOT EDIT!

[Unit]
Description=Dashing Dashboard
Before=systemd-user-sessions.service

[Service]
TimeoutStartSec=0
Environment=DASHBOARD_DIR=<%= @dashboard_path %>
Environment=DASHING_PORT=<%= @dashboard_port %>
WorkingDirectory=<%= @dashboard_path %>

ExecStart=<%= node['undef-dashing']['start_script'] %> start -d -e production -p <%= @dashboard_port %> -D

Type=oneshot
RemainAfterExit=yes
User=<%= node['undef-dashing']['user'] %>

ExecStop=<%= node['undef-dashing']['start_script'] %> stop
User=<%= node['undef-dashing']['user'] %>
